# autorippr@.service
[Unit]
Description=AutoRippr

[Service]
Type=oneshot
ExecStart=python /opt/autorippr/autorippr.py --all

[Install]
WantedBy=multi-user.target