# autorippr@.service
[Unit]
Description=AutoRippr

[Service]
Type=oneshot
ExecStart="$AUTORIPPR_HOME"/autorippr_wrapper.sh

[Install]
WantedBy=multi-user.target